# cfgs/etc/cron.daily/try-tarantool-org.cron
#!/bin/sh

/usr/local/try-tarantool-org/container/tool.sh cron
