# cfgs/etc/cron.daily/try_tarantool.cron
#!/bin/sh

/usr/local/tarantool/try_tarantool/container/try_tarantool_container.sh cron
